# .config/systemd/user/sync-org.service
[Service]
ExecStart=/bin/sh -c ${HOME}/bin/sync-org
Restart=always
NoNewPrivileges=false
